# socket-activation/templates/cache-EXAMPLE.socket
[Unit]
Description=Socket for EXAMPLE cache


[Socket]

ListenStream=GLOBAL_IP:PORT_EXAMPLE
Transparent=True

[Install]
WantedBy=sockets.target